# files/rc.local
iw dev wlan0 interface add uap0 type __ap
systemctl start hostapd
